# external/overlay/etc/s6/nginx/run
#!/bin/sh

mkdir -p /tmp/nginx
rm -f /var/cache/nginx
mkdir -p /var/cache
ln -s /var/cache/nginx /tmp/nginx

exec /sbin/nginx -c /etc/nginx/nginx.conf
