# external/overlay/etc/s6/tui-server/run
#!/bin/sh

mkdir -p /root/data
rm -f /tmp/tui-server.sock
exec /bin/tui-server -d /root/data/tui.db -u /tmp/tui-server.sock
